# mt5-bridge/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ── Virtual framebuffer (always on) ──────────────────────────────────
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1024x768x16
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── VNC server (setup mode only) ─────────────────────────────────────
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -nopw -listen 0.0.0.0 -rfbport 5900
autorestart=true
priority=20
autostart=%(ENV_ENABLE_VNC)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── noVNC web proxy (setup mode only) ────────────────────────────────
[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=25
autostart=%(ENV_ENABLE_VNC)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── MT5 terminal ─────────────────────────────────────────────────────
# Starts after MT5 is installed. On first run, start manually via noVNC.
# After installation, supervisord auto-restarts the terminal.
[program:mt5]
command=wine "C:/Program Files/MetaTrader 5/terminal64.exe" /portable
autorestart=true
priority=30
startsecs=10
startretries=5
environment=DISPLAY=":99",WINEPREFIX="/root/.wine"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── RPyC bridge server ──────────────────────────────────────────────
# Exposes MT5 Python API over RPyC for the telebot to connect to.
[program:rpyc]
command=wine C:/Python312/python.exe -m mt5linux --host 0.0.0.0 -p %(ENV_RPYC_PORT)s
autorestart=true
priority=40
startsecs=15
startretries=10
environment=DISPLAY=":99",WINEPREFIX="/root/.wine"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
